# debos-recipes/overlays/arduino-unoq-releases/etc/systemd/system/key_events.service
[Unit]
Description=Provides bindings for volumeup and volumedown keys on boot
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/arduino/scripts/keybinding_check.sh

[Install]
WantedBy=multi-user.target
