# debos-recipes/overlays/arduino-unoq-releases/etc/udev/rules.d/10-leds.rules
SUBSYSTEM=="leds", ACTION=="add", RUN+="/bin/chgrp -R gpiod /sys%p", RUN+="/bin/chmod -R g=u /sys%p"
SUBSYSTEM=="leds", ACTION=="change", RUN+="/bin/chgrp -R gpiod /sys%p", RUN+="/bin/chmod -R g=u /sys%p"
